# files/mirror.cron
0 */4 * * * nginx /usr/local/bin/mirror.py 2>&1 | logger -t "mirror" -p daemon.info
